# docker/rootfs/etc/services.d/couchdb/run
#!/usr/bin/with-contenv bash

echo "starting couchdb"
/docker-entrypoint.sh /opt/couchdb/bin/couchdb
